# etc/systemd/ocre-rdv-reminder.service
[Unit]
Description=Ocre Immo RDV reminder PWA push (M89)
After=network.target mariadb.service

[Service]
Type=oneshot
User=www-ocre
Group=www-data
WorkingDirectory=/opt/ocre-app/api
ExecStart=/usr/bin/php /opt/ocre-app/api/cron_rdv_reminder.php
StandardOutput=append:/var/log/ocre-rdv-reminder.log
StandardError=append:/var/log/ocre-rdv-reminder.log
